# 后端服务/captcha-server.service
[Unit]
Description=Local Captcha Recognition Service
Documentation=https://github.com/your-project
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=captcha
Group=captcha
WorkingDirectory=/opt/captcha-server
Environment=PYTHONPATH=/opt/captcha-server
Environment=FLASK_ENV=production
ExecStart=/opt/captcha-server/venv/bin/python start_production.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
RestartSec=5
Restart=always

# 安全配置
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
# 允许写入的路径（配置文件、规则文件、API Keys、日志文件）
ReadWritePaths=/opt/captcha-server
ReadWritePaths=/opt/captcha-server/*.json
ReadWritePaths=/opt/captcha-server/*.log
NoNewPrivileges=true

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=captcha-server

[Install]
WantedBy=multi-user.target